# fatmo.service
[Unit]
Description=Fatmo

[Service]
TimeoutStartSec=0
ExecStart=fatmo.py

[Install]
WantedBy=multi-user.target

